# eye.timer
[Unit]
Description=Run eye.service after 20 minutes

[Timer]
OnActiveSec=1200

[Install]
WantedBy=timers.target

