# eye.service
[Unit]
Description=Announce the end of timer in new urxvtc

[Service]
Type=oneshot
ExecStart=/usr/bin/urxvtc -e restart_eye_timer.sh
