# eye.timer
[Unit]
Description=Run eye.service every 20.33 minutes

[Timer]
OnActiveSec=12
OnUnitActiveSec=1220

[Install]
WantedBy=timers.target

