# eye.timer
[Unit]
Description=Run eye.service after 20 minutes

[Timer]
OnActiveSec=20 minutes

[Install]
WantedBy=timers.target

